# buildroot/br2_hudos/package/hudOS-first-boot/first_boot.service
[Unit]
Description=First Boot
Before=default.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/first_boot.sh
StandardOutput=tty

[Install]
WantedBy=default.target