# apps/japokwm/hud-apps.target
[Unit]
Description=hudOS Apps
After=hud.target
After=japokwm.target

[Install]
WantedBy=hud.target