# buildroot/package/hudOS/nreal-glasses.target
[Unit]
Description=Nreal glasses
StopWhenUnneeded=true
