# apps/weston/hud-apps.target
[Unit]
Description=hudOS Apps

[Install]
WantedBy=hud.target