# apps/weston/hud-apps.target
[Unit]
Description=hudOS Apps
After=hud.target

[Install]
WantedBy=hud.target