# apps/weston/hud-apps.target
[Unit]
Description=hudOS Apps
After=hud.target
After=weston.target

[Install]
WantedBy=hud.target