# apps/weston/hud.target
[Unit]
Description=hudOS

[Install]
WantedBy=default.target