# apps/weston/hud.target
[Unit]
Description=hudOS
After=multi-user.target

[Install]
WantedBy=multi-user.target