# apps/kiwmi/hud-apps.target
[Unit]
Description=hudOS Apps
After=hud.target
After=kiwmi.target

[Install]
WantedBy=hud.target