# apps/kiwmi/hud-apps.target
[Unit]
Description=hudOS Phone screen apps
After=hud.target

[Install]
WantedBy=hud.target