# apps/kiwmi/hud-glasses.target
[Unit]
Description=hudOS Glasses apps
After=hud.target
After=kiwmi.target
StopWhenUnneeded=true

[Install]
WantedBy=nreal-glasses.target