# apps/kiwmi/hud.slice
[Unit]
Description=hudOS Slice
Before=slices.target
